# systemd/h2v2-api.service
[Unit]
Description=H2V2 API
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=h2v2
Group=h2v2
EnvironmentFile=/opt/h2v2/.env.generated
RuntimeDirectory=h2v2
RuntimeDirectoryMode=0750
WorkingDirectory=/opt/h2v2/current
ExecStart=/opt/h2v2/bin/panel-api serve
Restart=always
RestartSec=2

[Install]
WantedBy=multi-user.target
